# 08-logging/08-files/agama-client.service
[Unit]
Description=<FIXME>
Documentation=<FIXME>
After=network-online.target

[Service]
User=<FIXME>
ExecStart=/usr/local/bin/<FIXME>

[Install]
WantedBy=multi-user.target
